# deploy/windows/intraday_1230.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  SEPA Signals — intraday scan at 17:30 local (mid-session).
  Mid-session check: catches setups that develop after the morning.

  StartBoundary time is in the machine's local clock — Task Scheduler does not
  convert timezones. The mini PC's clock is set to UK local time (BST/GMT).
  BST: 17:30 = 12:30 EDT (NYSE mid-session).

  Import with:
    schtasks /create /tn "SEPA-Intraday-1230" /xml intraday_1230.xml /f

  The <WorkingDirectory> path must match where sepa-signals is installed.
  Edit if your install path differs from C:\sepa-signals.
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>SEPA Signals intraday breakout scan — 17:30 BST (12:30 EDT market local time)</Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-01-01T17:30:00</StartBoundary>
      <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\lathe\sepa-signals\.venv\Scripts\python.exe</Command>
      <Arguments>-m sepa.run_intraday</Arguments>
      <WorkingDirectory>C:\Users\lathe\sepa-signals</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
