# deploy/windows/intraday_0945.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  SEPA Signals — intraday scan at 09:45 ET (14:45 London).
  Runs 15 minutes after the open so the first 5-minute bars have settled.

  Import with:
    schtasks /create /tn "SEPA-Intraday-0945" /xml intraday_0945.xml /f

  The <WorkingDirectory> path must match where sepa-signals is installed.
  Edit if your install path differs from C:\sepa-signals.
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>SEPA Signals intraday breakout scan — 09:45 ET (14:45 London)</Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-01-01T14:45:00</StartBoundary>
      <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\lathe\sepa-signals\.venv\Scripts\python.exe</Command>
      <Arguments>-m sepa.run_intraday --mode intraday</Arguments>
      <WorkingDirectory>C:\Users\lathe\sepa-signals</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
